# services/startup.service
[Unit]
Description=Hacky startup

[Service]
Type=oneshot
ExecStart=/usr/local/bin/startup.py

[Install]
WantedBy=multi-user.target
